# run/supervisor.ini
[program:LicensesPlates]
command=/home/pi/venv/venv1/bin/gunicorn_django -c /home/pi/Documents/projects/django-alpr/run/gunicorn_conf.py
directory=/home/pi/Documents/projects/django-alpr
user=pi
autostart=true
autorestart=true
priority=991
stopsignal=KILL